# local-agents/supervisord.conf
[unix_http_server]
file=/Users/jimmymalhan/stacky-os/logs/supervisor.sock

[supervisord]
logfile=/Users/jimmymalhan/stacky-os/logs/supervisord.log
logfile_maxbytes=20MB
logfile_backups=3
loglevel=info
pidfile=/Users/jimmymalhan/stacky-os/logs/supervisord.pid
nodaemon=false
user=jimmymalhan

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///Users/jimmymalhan/stacky-os/logs/supervisor.sock

[group:core]
programs=stacky-api

[group:agents]
programs=stacky-supervisor,stacky-agent-1,stacky-agent-2

[group:ui]
programs=stacky-frontend

[program:stacky-api]
command=/Library/Developer/CommandLineTools/usr/bin/python3 -m uvicorn api:app --host 0.0.0.0 --port 8000 --loop asyncio
directory=/Users/jimmymalhan/stacky-os
environment=BOS_HOME="/Users/jimmymalhan/stacky-os",CODEBASE_PATH="/Users/jimmymalhan/Documents/stacky",LOCAL_LIGHT="qwen2.5-coder:7b",LOCAL_STANDARD="llama3.1:8b",PORT_API="8000",OLLAMA_API_BASE="http://127.0.0.1:11434",PATH="/Users/jimmymalhan/Library/Python/3.9/bin:/usr/local/bin:/usr/bin:/bin"
autorestart=true
startretries=10
startsecs=3
redirect_stderr=true
stdout_logfile=/Users/jimmymalhan/stacky-os/logs/api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

[program:stacky-supervisor]
command=/Library/Developer/CommandLineTools/usr/bin/python3 /Users/jimmymalhan/stacky-os/supervisor-agent.py
directory=/Users/jimmymalhan/stacky-os
environment=BOS_HOME="/Users/jimmymalhan/stacky-os",PROJECT_SLUG="stacky",PROJECT_NAME="Stacky",PORT_API="8000",LOCAL_LIGHT="qwen2.5-coder:7b",LOCAL_STANDARD="llama3.1:8b",CODEBASE_PATH="/Users/jimmymalhan/Documents/stacky",OLLAMA_API_BASE="http://127.0.0.1:11434",PATH="/Users/jimmymalhan/Library/Python/3.9/bin:/usr/local/bin:/usr/bin:/bin"
autorestart=true
startretries=5
redirect_stderr=true
stdout_logfile=/Users/jimmymalhan/stacky-os/logs/supervisor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

[program:stacky-agent-1]
command=/bin/bash /Users/jimmymalhan/stacky-os/agent-loop.sh
directory=/Users/jimmymalhan/stacky-os
environment=BOS_HOME="/Users/jimmymalhan/stacky-os",WORKER_ID="1",PORT_API="8000",PROJECT_SLUG="stacky",PROJECT_NAME="Stacky",CODEBASE_PATH="/Users/jimmymalhan/Documents/stacky",LOCAL_LIGHT="qwen2.5-coder:7b",LOCAL_STANDARD="llama3.1:8b",OLLAMA_API_BASE="http://127.0.0.1:11434",PATH="/Users/jimmymalhan/Library/Python/3.9/bin:/usr/local/bin:/usr/bin:/bin"
autorestart=unexpected
startretries=3
redirect_stderr=true
stdout_logfile=/Users/jimmymalhan/stacky-os/logs/agent.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5

[program:stacky-agent-2]
command=/bin/bash /Users/jimmymalhan/stacky-os/agent-loop.sh
directory=/Users/jimmymalhan/stacky-os
environment=BOS_HOME="/Users/jimmymalhan/stacky-os",WORKER_ID="2",PORT_API="8000",PROJECT_SLUG="stacky",PROJECT_NAME="Stacky",CODEBASE_PATH="/Users/jimmymalhan/Documents/stacky",LOCAL_LIGHT="qwen2.5-coder:7b",LOCAL_STANDARD="llama3.1:8b",OLLAMA_API_BASE="http://127.0.0.1:11434",PATH="/Users/jimmymalhan/Library/Python/3.9/bin:/usr/local/bin:/usr/bin:/bin"
autorestart=unexpected
startretries=3
redirect_stderr=true
stdout_logfile=/Users/jimmymalhan/stacky-os/logs/agent.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5

[program:stacky-frontend]
command=/Library/Developer/CommandLineTools/usr/bin/python3 -m http.server 3000 --directory /Users/jimmymalhan/stacky-os/frontend
directory=/Users/jimmymalhan/stacky-os
autorestart=true
redirect_stderr=true
stdout_logfile=/Users/jimmymalhan/stacky-os/logs/frontend.log
stdout_logfile_maxbytes=5MB
